# install/asda-client.service
[Unit]
Description=ASDA Client
After=network.target fail2ban.service
Wants=fail2ban.service

[Service]
Type=simple
User=asda
Group=asda
WorkingDirectory=/opt/asda/client
ExecStart=/usr/bin/node /opt/asda/client/client.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=asda-client
Environment=NODE_ENV=production
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HOME=/home/asda
# Ensure sudo preserves environment for the asda user
Environment=SUDO_USER=asda

[Install]
WantedBy=multi-user.target
